# supervisord.conf
[supervisord]
nodaemon=true

[program:share]
command=python /app/share.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
